# systemd/nuc/ai-gateway.service
[Unit]
Description=AI Gateway
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=4

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/bin/bash -c "cd /opt/ai-gateway && /home/ubuntu/.local/bin/uv run main.py"

[Install]
WantedBy=multi-user.target
